# seafile/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -euo pipefail

for _ in $(seq 1 60); do
    if [ -f /etc/nginx/sites-available/seafile.conf ]; then
        break
    fi
    sleep 1
done

bashio::log.info "Starting Nginx..."
nginx -t
exec nginx -g 'daemon off;'
